# sys_sensors_mqtt.service
[Unit]
Description=SysSensorsMQTT

[Service]
Type=idle
WorkingDirectory=/home/osmc/SysSensorsMQTT
ExecStart=/home/osmc/SysSensorsMQTT/sys_sensors_mqtt_daemon.py

[Install]
WantedBy=multi-user.target
